# systemd/survey-api.socket
[Unit]
Description=gunicorn socket survey-api

[Socket]
ListenStream=/run/survey-api.sock

[Install]
WantedBy=sockets.target
